# ias-setup.service
[Unit]
Description=ias-setup script one-time

[Service]
Type=oneshot
ExecStart=/usr/bin/ias-setup

[Install]
WantedBy=basic.target