# ias-setup.service
[Unit]
Description=ias-setup script one-time
Before=ias.service

[Service]
ExecStart=/usr/bin/ias-setup

[Install]
WantedBy=ias.service